# ops/systemd/cabinet.service.tmpl
[Unit]
Description=Cabinet AI Workspace
After=basic.target
StartLimitIntervalSec=60
StartLimitBurst=10

[Service]
Type=simple

WorkingDirectory=@CABINET_ROOT@

Environment=HOME=@HOME@
Environment=XDG_CONFIG_HOME=@HOME@/.config
Environment=XDG_STATE_HOME=@HOME@/.local/state

Environment=NPM_CONFIG_PREFIX=@HOME@/.npm-global
Environment=PATH=@HOME@/.npm-global/bin:@HOME@/.local/bin:@HOME@/bin:/usr/local/bin:/usr/bin:/bin
Environment=ADAPTER_RUNTIME_PATH=@HOME@/.npm-global/bin:@HOME@/.local/bin:@HOME@/bin:/usr/local/bin:/usr/bin:/bin

Environment=NEXT_TELEMETRY_DISABLED=1
Environment=CABINET_TELEMETRY_DISABLED=1

Environment=CABINET_APP_PORT=4001
Environment=CABINET_DAEMON_PORT=4100

EnvironmentFile=-@HOME@/.config/cabinet/runtime.env

ExecStartPre=/usr/bin/test -x @HOME@/.local/bin/cabinet-session
ExecStart=@HOME@/.local/bin/cabinet-session

Restart=always
RestartSec=5

TimeoutStartSec=90
TimeoutStopSec=30
KillMode=control-group

LimitNOFILE=65536
UMask=0077

[Install]
WantedBy=default.target
