# ops/systemd/systemkatalog.service.tmpl
[Unit]
Description=Systemkatalog
After=basic.target
StartLimitIntervalSec=60
StartLimitBurst=5

[Service]
Type=simple
WorkingDirectory=@SYSTEMKATALOG_ROOT@
Environment=HOME=@HOME@
Environment=PYTHONDONTWRITEBYTECODE=1
Environment=SYSTEMKATALOG_ROOT=@SYSTEMKATALOG_ROOT@
Environment=SYSTEMKATALOG_BIND=127.0.0.1
Environment=SYSTEMKATALOG_PORT=4001
ExecStartPre=/usr/bin/python3 @SYSTEMKATALOG_ROOT@/scripts/validate_system_catalog.py
ExecStartPre=/usr/bin/python3 @SYSTEMKATALOG_ROOT@/scripts/render_system_catalog.py --check
ExecStartPre=/usr/bin/python3 @SYSTEMKATALOG_ROOT@/scripts/serve_system_catalog.py --check
ExecStart=@HOME@/.local/bin/systemkatalog
Restart=on-failure
RestartSec=3
TimeoutStartSec=30
TimeoutStopSec=15
KillMode=control-group
UMask=0077
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
LockPersonality=true

[Install]
WantedBy=default.target
