# ops/systemd/cabinet.service.d/10-loopback-gate.conf.tmpl
[Service]
ExecStartPre=@HOME@/.local/bin/cabinet-security-gate
